# met-cron/cron/crontab
* * * * * default cd /met-cron && ./run_met_etl.sh
0 8 * * * default cd /met-cron && ./run_met_closeout.sh
# An empty line is required at the end of this file for a valid cron file