# met-cron/cron/crontab
# ANALYTICS ETL Runs At every hour (UTC).
0 * * * * default cd /met-cron && ./run_met_etl.sh
# ENGAGEMENT CLOSEOUT Runs At 17:00 on every day-of-week from Monday through Friday (UTC).
0 17 * * 1-5 default cd /met-cron && ./run_met_closeout.sh
# ENGAGEMENT PUBLISH Runs At every 5 minutes.
*/5 * * * * default cd /met-cron && ./run_met_publish.sh
# PURGE Runs At every 15 days.
0 0 * * 0 default cd /met-cron && ./run_met_purge.sh
# REDACT COMMENTS Runs At every day.
0 0 */1 * * default cd /met-cron && ./run_met_comment_redact.sh
# An empty line is required at the end of this file for a valid cron file